# www-base/supervisord-ssh.conf
[program:ssh]
command=/usr/sbin/sshd -D -p 2222
numprocs=1
autostart=true
autorestart=true
priority=1
