# loadbalancer/supervisord-nginx.conf
[program:nginx]
command=/start-nginx.sh
